# spotify/rootfs/etc/services.d/spotifyd/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Spotify Connect
# Runs librespot
# ==============================================================================
declare -a options
declare bitrate
declare name
declare password
declare username
declare cache_size_limit
declare cache_path
declare disable_audio_cache
declare device_type
declare format
declare initial_volume
declare enable_volume_normalisation
declare autoplay

bashio::log.info 'Starting the Spotify daemon...'

# Bitrate
bitrate=$(bashio::config 'bitrate')
options+=(--bitrate "${bitrate}")

# Device name
name=$(bashio::config 'name')
options+=(--name "${name}")

# Username / password
if bashio::config.has_value 'username'; then
    bashio::config.require.password
    username=$(bashio::config 'username')
    password=$(bashio::config 'password')
    options+=(--username "${username}")
    options+=(--password "${password}")
fi

# Cache path and size limit
if bashio::config.has_value 'cache-size-limit'; then
    cache_size_limit=$(bashio::config 'cache-size-limit')
    if [[ "$cache_size_limit" =~ ^[0-9]+[KMG]?$ ]]; then
        options+=(--cache-size-limit "${cache_size_limit}")
    else
        bashio::log.warning "Invalid cache-size-limit, using default of 1M"
        options+=(--cache-size-limit "1M")
    fi
fi

# Set cache path to enable cache-size-limit
cache_path="/tmp/librespot_cache"
options+=(--cache "${cache_path}")

# Disable audio cache
if bashio::config.true 'disable-audio-cache'; then
    options+=(--disable-audio-cache)
fi

# Device type
device_type=$(bashio::config 'device-type')
if [[ -z "${device_type}" ]]; then
    device_type="speaker" # Default to "speaker" if not provided
fi
options+=(--device-type "${device_type}")

# Format
format=$(bashio::config 'format')
options+=(--format "${format}")

# Initial volume
initial_volume=$(bashio::config 'initial-volume')
if [[ "$initial_volume" =~ ^[0-9]+$ ]]; then
    options+=(--enable-volume-normalisation)
    options+=(--initial-volume "${initial_volume}")
else
    bashio::log.warning "Invalid initial-volume, using default of 50"
    options+=(--enable-volume-normalisation)
    options+=(--initial-volume "50")
fi

# Autoplay
if bashio::config.true 'autoplay'; then
    options+=(--autoplay)
fi

# Are we running in debug mode?
if bashio::debug; then
    options+=(--verbose)
fi

# Log the final options for debugging purposes
bashio::log.info "Running librespot with options: ${options[*]}"

# Run librespot
exec librespot "${options[@]}"
